# g-nginx-dynamic-whitelists/g-nginx-dynamic-whitelists-1.0/debian/g-nginx-dynamic-whitelists.cron.d
*/15 * * * * root /usr/bin/g-nginx-dynamic-whitelists
